# deploy/systemd/seevar-next-dashboard.service
[Unit]
Description=SeeVar Next dashboard

[Service]
WorkingDirectory=%h/SeeVar-Next
ExecStart=%h/SeeVar-Next/.venv/bin/seevar-next-dashboard
Restart=on-failure

[Install]
WantedBy=default.target
