# deploy/systemd/seevar-next-preflight.service
[Unit]
Description=SeeVar Next nightly preflight

[Service]
Type=oneshot
WorkingDirectory=%h/SeeVar-Next
ExecStart=%h/SeeVar-Next/.venv/bin/seevar-next preflight
