# deploy/systemd/seevar-next-nightly.timer
[Unit]
Description=Run SeeVar Next automated nightly run

[Timer]
OnCalendar=*-*-* 22:30:00
Persistent=false

[Install]
WantedBy=timers.target
