# deploy/systemd/seevar-next-nightly.service
[Unit]
Description=SeeVar Next automated nightly run

[Service]
Type=oneshot
WorkingDirectory=%h/SeeVar-Next
ExecStart=%h/SeeVar-Next/.venv/bin/seevar-next preflight
ExecStart=%h/SeeVar-Next/.venv/bin/seevar-next readiness
ExecStart=%h/SeeVar-Next/.venv/bin/seevar-next-flight policy
ExecStart=%h/SeeVar-Next/.venv/bin/seevar-next-flight steps --plan data/tonights_plan.json --proof data/flight_runs/flight_steps.jsonl
ExecStart=%h/SeeVar-Next/.venv/bin/seevar-next-flight submit --plan data/seestarpy_plan.json
